# supervisord-elasticsearch.conf
[program:elasticsearch]
command=/start-es.sh
# environment=USER="elsearch"
# directory=/opt/elasticsearch-6.6.2/
numprocs=1
autostart=true
autorestart=false
user=root